# vamp/vagrant/modules/marathon/templates/marathon.service.erb
[Unit]
Description=Marathon
After=network.target
Wants=network.target

<% @values = { 'marathon::marathon_assets_path' => scope['marathon::marathon_assets_path'],
               'marathon::marathon_checkpoint' =>  scope['marathon::marathon_checkpoint'],
               'marathon::marathon_event_subscribe' => scope['marathon::marathon_event_subscriber'],
               'marathon::marathon_executor' => scope['marathon::marathon_executor'],
               'marathon::marathon_failover_timeout' => scope['marathon::marathon_failover_timeout'],
               'marathon::marathon_ha' => scope['marathon::marathon_ha'],
               'marathon::marathon_hostname' => scope['marathon::marathon_hostname'],
               'marathon::marathon_http_credentials' => scope['marathon::marathon_http_credentials'],
               'marathon::marathon_http_endpoints' => scope['marathon::marathon_http_endpoints'],
               'marathon::marathon_http_port' => scope['marathon::marathon_http_port'],
               'marathon::marathon_https_port' => scope['marathon::marathon_https_port'],
               'marathon::marathon_local_port_max' => scope['marathon::marathon_local_port_max'],
               'marathon::marathon_local_port_min' => scope['marathon::marathon_local_port_min'],
               'marathon::marathon_master' => scope['marathon::marathon_master'],
               'marathon::marathon_mesos_role' => scope['marathon::marathon_mesos_role'],
               'marathon::marathon_mesos_user' => scope['marathon::marathon_mesos_user'],
               'marathon::marathon_reconciliation_frequency' => scope['marathon::marathon_reconciliation_frequency'],
               'marathon::marathon_reconciliation_initial_delay' => scope['marathon::marathon_reconciliation_initial_delay'],
               'marathon::marathon_ssl_keystore_password' => scope['marathon::marathon_ssl_keystore_password'],
               'marathon::marathon_ssl_keystore_path' => scope['marathon::marathon_ssl_keystore_path'],
               'marathon::marathon_task_launch_timeout' => scope['marathon::marathon_task_launch_timeout'],
               'marathon::marathon_task_rate_limit' => scope['marathon::marathon_task_rate_limit'],
               'marathon::marathon_zk' => scope['marathon::marathon_zk'],
               'marathon::marathon_zk_hosts' => scope['marathon::marathon_zk_hosts'],
               'marathon::marathon_zk_state' => scope['marathon::marathon_zk_state'],
               'marathon::marathon_zk_timeout' => scope['marathon::marathon_zk_timeout'],
             } -%>
[Service]
ExecStart=<%= @install_dir %>/bin/start <% @values.each do |key,val| -%>
<% if val != :undef -%>
 --<% arg = key.sub('marathon::marathon_', '') -%><%= arg -%> <%= val -%>
<% end -%>
<% end %>
Restart=on-abort
Restart=always
RestartSec=20

[Install]
WantedBy=multi-user.target
